# .deploy/supervisor/supervisor-local.ini
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:events]
process_name = events-%(process_num)s
command=go run /Volumes/Source/comercial/lpgenerator/ws-notify/events.go
directory=/Volumes/Source/comercial/lpgenerator/ws-notify/
autostart=true
autorestart=true
numprocs = 1
stopsignal=QUIT
stopwaitsecs=360
stopasgroup = true
redirect_stderr=True
stdout_logfile=/Volumes/Source/comercial/lpgenerator/ws-notify/lpg.events.log
